# configs/stage3_ubuntu/systemd/getty@tty1.service.d/10_overrides.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin mlab --noclear %I $TERM

