# configs/stage3_ubuntu/systemd/enable-bbr.service
[Unit]
Description=Enable tcp_bbr module.

[Service]
Type=oneshot
ExecStart=/usr/sbin/modprobe tcp_bbr

