# configs/stage3_ubuntu/systemd/generate-eth0-config.service
[Service]
Type=oneshot
ExecStart=/usr/share/oem/generate_network_config.sh /etc/systemd/network/00-eth0.network
