# configs/stage3_ubuntu/systemd/mlab-periodic-job.service
[Unit]
Description=Runs periodic jobs.
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/bin/run-parts /etc/periodic/15min

