# configs/stage2/inittab
# First initscript to run after init starts.
::sysinit:/etc/rc.local

# TTYs (the first is default console)
tty1::askfirst:/bin/sh
tty2::respawn:/bin/login
tty3::respawn:/bin/login
tty4::respawn:/bin/login

# Commands to run when restarting init.
::restart:/sbin/init

# Commands to run before shutdown or reboot.
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
