# configs/virtual_ubuntu/etc/systemd/system/join-cluster.service
[Unit]
Description=Join the M-Lab Platform Cluster
# The presence of this file indicates that this machine has already joined the
# cluster.  If it exists, do not run this unit.
ConditionPathExists=!/etc/kubernetes/kubelet.conf

[Service]
Type=oneshot
ExecStart=/opt/mlab/bin/join-cluster.sh

[Install]
WantedBy=multi-user.target
