# configs/virtual_ubuntu/etc/systemd/system/reboot-api-node.service
[Unit]
Description=reboot-api-node.service

[Service]
Type=oneshot
ExecStart=/opt/mlab/bin/reboot-api-node.sh

