# configs/stage3_ubuntu/etc/systemd/system/check-reboot.service
[Unit]
Description=Check system uptime, and reboot the machine if it is too long

[Service]
Type=oneshot
ExecStart=/opt/mlab/bin/check-reboot.sh

[Install]
WantedBy=multi-user.target

