# configs/stage3_ubuntu/etc/systemd/system/mlab-set-quotas.timer
[Unit]
Description=Runs mlab-set-quotas.service every 15 minutes.

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=multi-user.target

