# configs/stage3_ubuntu/etc/systemd/system/containerd.service.d/10-overrides.conf
[Unit]
# The containerd service depends on the /cache/containerd volume.
Requires=cache-containerd.mount
After=cache-containerd.mount

