# configs/stage3_ubuntu/etc/systemd/system/containerd.service.d/20-limit-nofile.conf
# containerd 2.x (which arrived on platform machines via the Ubuntu jammy
# security update that bumped the containerd package from 1.7.x to 2.x) no
# longer sets a high LimitNOFILE on its packaged systemd unit. As a result the
# containerd daemon — and every container that inherits the daemon's rlimits —
# defaults to a soft RLIMIT_NOFILE of only 1024. Services that do not raise
# their own soft limit at startup (e.g. pusher) then exhaust it, which silently
# stalls data uploads and loses data. Restore the previously-effective high
# default on the daemon so all containers inherit it again.
#
# Note: this must set a high *soft* limit (not just hard), because the affected
# services do not raise their own soft limit. LimitNOFILE=N sets soft=hard=N.
[Service]
LimitNOFILE=1048576
